# containers/rootfs/etc/s6-overlay/s6-rc.d/vaultwarden/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck -d

cd /vw/
/usr/bin/exec /vw/vaultwarden