# containers/rootfs/etc/cron.d/backup
01 21 * * * root /command/with-contenv /usr/local/sbin/backup >/proc/1/fd/1 2>/proc/1/fd/2
